# data/com.redhat.Yggdrasil1.Worker1.foreman.service.in
[Unit]
Description=Foreman remote execution yggdrasil worker service
Documentation=https://github.com/theforeman/foreman_ygg_worker

[Service]
Type=dbus
User=root
Group=@worker_group@
ExecStart=@libexecdir@/foreman_worker
BusName=com.redhat.Yggdrasil1.Worker1.foreman

[Install]
WantedBy=multi-user.target
